# spicerack/sesame/cpp/etc/sesame.service
[Unit]
Description=MRG Management system agent

[Service]
Type=simple
ExecStart=/usr/bin/sesame

[Install]
WantedBy=multi-user.target
